# macos/plists/com.xray-cfg.xray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xray-cfg.xray</string>
    <key>UserName</key>
    <string>root</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/macos/runtime/bin/xray/xray</string>
        <string>run</string>
        <string>-c</string>
        <string>__REPO_ROOT__/macos/runtime/bin/xray/config.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__/macos/runtime/bin/xray</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>XRAY_LOCATION_ASSET</key>
        <string>__REPO_ROOT__/macos/runtime/bin/xray</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/macos/runtime/xray-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/macos/runtime/xray-stderr.log</string>
</dict>
</plist>
